# systemd/chalna-web.service
[Unit]
Description=Chalna Web UI (Gradio)
After=network.target

[Service]
Type=simple
User=jonhpark
WorkingDirectory=/home/jonhpark/workspace/chalna
Environment="PATH=/home/jonhpark/workspace/chalna/venv/bin"
ExecStart=/home/jonhpark/workspace/chalna/venv/bin/chalna web --host 0.0.0.0 --port 7860
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
